# examples/20231008-automaticbackup/etc/udev/rules.d/99-my-usb-rule.rules
KERNEL=="sd*", SUBSYSTEMS=="usb", ATTRS{idVendor}=="058f", ATTRS{idProduct}=="6387", RUN+="/root/open-close-backup-usb.sh"

